# etc/systemd/system/wifi-manager.service
[Unit]
Description=WiFi Manager - Timed WiFi Control
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wifi-manager.sh
StandardOutput=null
StandardError=journal
