# etc/systemd/system/wifi-unblock.service
[Unit]
Description=Ensure WiFi starts enabled on boot
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
# Clear any rfkill blocks that may have persisted
ExecStartPre=/bin/rm -f /var/lib/systemd/rfkill/*
# Unblock all rfkill
ExecStart=/usr/sbin/rfkill unblock all
ExecStart=/usr/sbin/rfkill unblock wifi
ExecStart=/usr/sbin/rfkill unblock 0
ExecStart=/usr/sbin/rfkill unblock 1
# Ensure interface is up
ExecStartPost=/sbin/ip link set wlan0 up
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target

