# etc/systemd/system/wifi-unblock.service
[Unit]
Description=Force-enable Wi-Fi at boot
After=sysinit.target
Before=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rfkill unblock wifi
ExecStartPost=/usr/bin/env sh -c 'command -v nmcli >/dev/null && nmcli radio wifi on'
RemainAfterExit=yes
